# scripts/memory_monitor.service
[Unit]
Description=Bible Clock Memory Monitor (DEPRECATED - Use bible-clock-system-monitor.service)
After=network.target bible-clock.service
Wants=bible-clock.service

[Service]
Type=simple
User=admin
Group=admin
ExecStart=/home/admin/Bible-Clock-v4/scripts/memory_monitor.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
WorkingDirectory=/home/admin/Bible-Clock-v4

[Install]
WantedBy=multi-user.target